# deploy/study-bot.service
[Unit]
Description=Discord Study Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/Discord-Bot
Environment="PATH=/opt/Discord-Bot/.venv/bin"
ExecStart=/opt/Discord-Bot/.venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/study-bot.log
StandardError=append:/var/log/study-bot.log

[Install]
WantedBy=multi-user.target
